# crontab
# Tidio Product Sync - supercronic crontab
# Incremental sync: every 2 hours, skipping 02:00 (replaced by full sync)
# Full sync: 02:00 daily
0 0,4,6,8,10,12,14,16,18,20,22 * * * /usr/local/bin/python /app/app.py
0 2 * * * /usr/local/bin/python /app/app.py --full
